# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
loglevel=info

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

# ============================================================
# CHATBOT: démarre EN PREMIER (priorité haute = chiffre bas)
# startsecs=30 car sentence-transformers prend du temps à charger
# ============================================================
[program:chatbot]
command=python -m uvicorn app.main:app --host 127.0.0.1 --port 8001 --log-level info --timeout-keep-alive 300
directory=/chatbot
environment=ANONYMIZED_TELEMETRY="False",OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",OPENAI_MODEL="%(ENV_OPENAI_MODEL)s",TN360_API_URL="%(ENV_TN360_API_URL)s",TEMPERATURE="%(ENV_TEMPERATURE)s",CHROMA_PERSIST_DIR="%(ENV_CHROMA_PERSIST_DIR)s",EMBEDDING_MODEL="%(ENV_EMBEDDING_MODEL)s"
autostart=true
autorestart=true
startretries=5
startsecs=30
stopwaitsecs=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

# ============================================================
# NGINX: démarre APRÈS le chatbot (priorité basse = chiffre haut)
# ============================================================
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
